# packages/taskqtool/taskqtool-0.0.4-py3-none-any.whl/task/statics/taskServer.service
[Unit]
Description=Task Server
Documentation=http://zhcloud.top
After=network-online.target
Wants=network-online.target

[Service]
User=root
Group=root
ExecStart=/usr/bin/python xxx.py
ExecReload=
ExecStop=

LimitMEMLOCK=infinity
# total threads (user+kernel)
TasksMax=infinity
TasksAccounting=false


[Install]
WantedBy=multi-user.target
